# tests/supervisord.conf
[supervisord]
user=root
nodaemon=true

[program:mlflow_server]
directory=%(ENV_DEST)s
command=bash -c "if [ %(ENV_RUN_APP)s = 'true' ]; then %(ENV_PYTHON_PATH)s/python %(ENV_DEST)s/main.py --mlflow_server %(ENV_MLFLOW_SERVER)s server; fi"
stdout_logfile =%(ENV_DEST)s/mlflow_server.log
redirect_stderr=true

[program:flask_app]
directory=%(ENV_DEST)s
command=bash -c "if [ %(ENV_RUN_APP)s = 'true' ]; then %(ENV_PYTHON_PATH)s/python %(ENV_DEST)s/main.py --mlflow_server %(ENV_MLFLOW_SERVER)s app --workers %(ENV_WORKERS)s; fi"
stdout_logfile =%(ENV_DEST)s/pymlframework_app.log
redirect_stderr=true